# config/systemd/thruster.service.erb
[Unit]
Description=Deploy app with Thruster
After=network.target sidekiq.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/starch
Environment=RAILS_ENV=production
Environment=TLS_DOMAIN=<%= Rails.application.credentials.domain %>
ExecStart=/usr/local/bin/thrust bin/rails server
Restart=always

StandardOutput=append:/root/starch/log/thruster-output.log
StandardError=append:/root/starch/log/thruster-error.log

[Install]
WantedBy=multi-user.target